# myground@.service
[Unit]
Description=MyGround self-hosting platform (user %i)
After=network-online.target docker.service
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
User=%i
ExecStart=/usr/local/bin/myground start
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
